# deploy/docker/supervisord.conf
#********************************************************************
# 项目名称：gin-vue3-blog
# 文件名称：supervisord.conf
# 创建时间：2026-03-17 22:18:35
#
# 系统用户：wyl
# 作　　者：無以菱
# 联系邮箱：huangjing510@126.com
# 功能描述：supervisord 进程管理配置，启动 blog-backend 与 nginx
#********************************************************************

[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/var/run/supervisord.pid

[program:blog-backend]
directory=/app
command=/app/blog-backend
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
